# agent/greenops-agent.service
[Unit]
Description=GreenOps Energy Intelligence Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=greenops
Group=greenops
WorkingDirectory=/opt/greenops-agent
ExecStart=/opt/greenops-agent/venv/bin/python greenops_agent.py --server https://YOUR_SERVER_URL
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=greenops-agent
Environment=PYTHONUNBUFFERED=1

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/greenops/.greenops

[Install]
WantedBy=multi-user.target
